# ai-trends.service
[Unit]
Description=AI Trends Service
After=network.target

[Service]
Type=forking
User=ubuntu
WorkingDirectory=/home/ubuntu/ai-trends
ExecStart=/usr/local/bin/pm2 start ecosystem.config.js --daemon
ExecReload=/usr/local/bin/pm2 reload ecosystem.config.js
ExecStop=/usr/local/bin/pm2 delete ai-trends
Restart=always
RestartSec=10
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target